# .config/systemd/user/eclim.service
[Unit]
Description=Eclim Daemon

[Service]
ExecStart=/home/tsuro/opt/eclipse/eclimd

